# crontab
*/30 * * * * root cd /app && python main.py >> /var/log/cron.log 2>&1

